# deploy/systemd/catshop-web.service
[Unit]
Description=CatShop Web (Next.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/project/shop/catshop
EnvironmentFile=/etc/catshop/catshop-web.env
ExecStart=/usr/bin/env node node_modules/next/dist/bin/next start --hostname 127.0.0.1
Restart=always
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
SyslogIdentifier=catshop-web

[Install]
WantedBy=multi-user.target
